# deploy/my-greeting-web.service
[Unit]
Description = my greeting web app
After = syslog.target

[Service]
ExecStart = /var/my-greeting-web/my-greeting-web-1.0.0.jar
Restart = always
Type = simple
User = root
Group = root
SuccessExitStatus = 143